# Jarvis_Development/docker/supervisord.rag.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid

[program:qdrant]
command=/qdrant/qdrant
autostart=true
autorestart=true
stdout_logfile=/app/logs/qdrant.log
stderr_logfile=/app/logs/qdrant.log
environment=QDRANT__SERVICE__HTTP_PORT=6333,QDRANT__STORAGE__STORAGE_PATH=/qdrant/storage

[program:chromadb]
command=uvicorn chromadb.app:app --host 0.0.0.0 --port 8001
autostart=true
autorestart=true
stdout_logfile=/app/logs/chroma.log
stderr_logfile=/app/logs/chroma.log
environment=IS_PERSISTENT="1",PERSIST_DIRECTORY="/chroma/data",ANONYMIZED_TELEMETRY="False"
